# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python2 /flask/start.py

[program:mongo]
command=mongod --smallfiles

